# karma/files/bo_servicios_initd
#!/bin/bash

# chkconfig: 2345 99 50
# description: script de control de servicios de BO


##############################################################################
# boservicios configuration file
##############################################################################
# NOTE: This file is automatically generated by puppet
# Changes to this file will be overwritten periodically by puppet!
##############################################################################

case $1 in
	start)
		cd /opt/bo/BI_40_2/sap_bobj
		su bo -c ./db2startup.sh
		su bo -c ./startservers
		su bo -c ./tomcatstartup.sh
		exit 0
		;;
	stop)
		cd /opt/bo/BI_40_2/sap_bobj
                su bo -c ./tomcatshutdown.sh
                su bo -c ./stopservers
                su bo -c ./db2shutdown.sh
		exit 0
		;;
	status)
		cd /opt/bo/BI_40_2/sap_bobj
		su bo -c ./status.sh
		exit 0
		;;
	*)
		echo "Usage boservicios [start|stop|status]"
		exit 1
		;;
esac
