# karma/files/tomcat6_initd
#!/bin/bash

##############################################################################
# tomcat6_initd file
##############################################################################
# NOTE: This file is automatically generated by puppet
# Changes to this file will be overwritten periodically by puppet!
##############################################################################


# chkconfig: 2345 55 25
# description: Tomcat
case "$1" in
        start)
                echo "Starting Tomcat..."
		pgrep -u tomcat java && exit 0
                su - tomcat -c /opt/tomcat/bin/startup.sh
        ;;
        stop)
                echo "Stopping Tomcat..."
		pgrep -u tomcat java || exit 0
                su - tomcat -c /opt/tomcat/bin/shutdown.sh
        ;;
        status)
                echo "Checking Tomcat..."
                pgrep -u tomcat java
                if [ $? -eq 0 ] ; then
                 echo "Tomcat s running..."
                 exit 0
                else
                 echo "Tomcat not running..."
                 exit 1
                fi
        ;;

        restart)
                su - tomcat -c /opt/tomcat/bin/shutdown.sh
                echo "Restarting Tomcat" ;
		while [ "`pgrep java`" != "" ] ; do 
		sleep 5 ;
		done
                su - tomcat -c /opt/tomcat/bin/startup.sh
        ;;
        *)
                echo "Usage: $0 {start|stop|restart}"
                exit 1
        ;;
esac
