# jboss5/files/jboss5_initd
#!/bin/bash

##############################################################################
# jboss5_initd configuration file
##############################################################################
# NOTE: This file is automatically generated by puppet
# Changes to this file will be overwritten periodically by puppet!
##############################################################################

# chkconfig: 2345 55 25
# description: jboss5
LOG="/var/log/jboss.log"
case "$1" in
        start)
                echo "Starting Jboss..."
		pgrep -u jboss java && exit 0 
                su - jboss -c "/opt/jboss/bin/run.sh &" >>$LOG
		exit 0 
        ;;
        stop)
                echo "Stopping Jboss..."
		pgrep -u jboss java ||  exit 0
                su - jboss -c "/opt/jboss/bin/shutdown.sh -S" >>$LOG
		sleep 10  ; pkill -u jboss java
		exit 0 
        ;;
        status)
                echo "Checking Jboss..."
                pgrep -u jboss java
                if [ $? -eq 0 ] ; then
                 echo "Jboss s running..."
                 exit 0
                else
                 echo "Jboss not running..."
                 exit 1
                fi

        ;;
        restart)
                su - jboss -c "/opt/jboss/bin/shutdown.sh -S" >>$LOG
		sleep 10  ; pkill -u jboss java ; pkill -u jboss java
                echo "Restarting Jboss" ;
                su - jboss -c "/opt/jboss/bin/run.sh &" >>$LOG
        ;;
        *)
                echo "Usage: $0 {start|stop|restart|status}"
                exit 1
        ;;
esac
